# CI/alchemy-next-git-ci.timer
[Unit]
Description=Alchemy Viewer CI Timer

[Timer]
OnUnitActiveSec=3s
OnBootSec=30s

[Install]
WantedBy=timers.target